# linux/systemd/qubes-session.target
[Unit]
Description=Qubes graphical session
Wants=qubes-guid.service
