# linux/systemd/qubes-qrexec@.socket
[Unit]
PartOf=qubes-qube@%i.service

[Socket]
ListenStream=/run/qubes/qrexec.%i
