# linux/systemd/qubesd.socket
[Socket]
ListenStream=/run/qubesd.sock
ListenStream=/run/qubesd.internal.sock
ListenStream=/run/qubesd.misc.sock

[Install]
WantedBy=sockets.target
